# 89-rpi-midi-hub.rules
ACTION=="add", SUBSYSTEM=="usb", DRIVER=="usb", RUN+="/etc/rpi-midi-hub/midiconnect.sh"
